# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:face_svr]
directory=/data/apps
command=python /data/apps/svr.py
autorestart=true